# packaging/scripts/chiton.service.in
[Unit]
Description=Chiton NVR
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=@bindir@/@PKGNAME@ -s
User=@runuser@
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT
SendSIGKILL=yes
TimeoutStopSec=15
Restart=on-failure
RuntimeDirectory=@PKGNAME@

[Install]
WantedBy=multi-user.target
